# claude/scheduled/taylor-morning-briefing/com.taylor.daily-briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Unique identifier for this launch agent -->
    <key>Label</key>
    <string>com.taylor.daily-briefing</string>

    <!-- Script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/taylor/Claude/Scheduled/taylor-morning-briefing/run-briefing.sh</string>
    </array>

    <!-- Run at 7:30 AM every day (local machine time) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- Do NOT run immediately on load — wait for 7:30 AM -->
    <key>RunAtLoad</key>
    <false/>

    <!-- launchd-level stdout/stderr (script also logs internally) -->
    <key>StandardOutPath</key>
    <string>/Users/taylor/Claude/briefing-output/.logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/taylor/Claude/briefing-output/.logs/launchd-stderr.log</string>

    <!-- Ensure common binary paths are visible to launchd (no shell profile loaded) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/taylor</string>
    </dict>

</dict>
</plist>
